# simpleLCD.service
[Unit]
 Description=My LCD Service
 After=multi-user.target

 [Service]
 Type=idle
 User=root
 ExecStartPre=/bin/sleep 30
<<<<<<< HEAD
 ExecStart=/bin/sh -c '/usr/bin/python3 /home/$(USER)/DWIN_T5UIC1_LCD/run.py  > /tmp/lcd.log 2>&1'
=======
 ExecStart=/bin/sh -c '/usr/bin/env python3 /home/pi/DWIN_T5UIC1_LCD/run.py  > /tmp/lcd.log 2>&1'
>>>>>>> b5055af842392453051d0e676bf273954ab09e28
 Restart=always

 [Install]
 WantedBy=multi-user.target

